# templates/appium.service.j2
[Unit]
Description=Appium
Documentation=https://appium.io
After=network.target

[Service]
Type=simple
User={{ appium_user }}
ExecStart=/bin/bash /home/{{ appium_user }}/run-appium

[Install]
WantedBy=multi-user.target
